# deckster.service
[Unit]
Description=Deckster

[Service]
Type=simple
ExecStart=%h/.local/bin/deckster

[Install]
WantedBy=default.target